# deploy/supervisor/eclipse-chat.conf
; Eclipse Chat — supervisor program for Fastify backend.
; Установка (initial-setup.sh делает это автоматически):
;   sudo cp deploy/supervisor/eclipse-chat.conf /etc/supervisor/conf.d/eclipse-chat.conf
;   sudo supervisorctl reread && sudo supervisorctl update

[program:eclipse-chat-server]
command=/usr/bin/node /var/www/eclipse-chat/apps/server/dist/index.js
directory=/var/www/eclipse-chat/apps/server
user=www-data
autostart=true
autorestart=true
startsecs=10
startretries=5
stderr_logfile=/var/log/supervisor/eclipse-chat.err.log
stdout_logfile=/var/log/supervisor/eclipse-chat.out.log
environment=NODE_ENV="production"
stopwaitsecs=10

; .env лежит в apps/server/.env и читается через dotenv в самом приложении.
; Если хочется передать переменные напрямую (не через .env) — добавляй
; их через environment="NODE_ENV=production,FOO=bar" одной строкой.
